# certbot.cron
30 00 * * 0 /usr/local/bin/letsencrypt.sh
