# gridftp-le/image/supervisord.ini
[supervisord]
nodaemon=true

[program:certbot-once]
command=/usr/local/bin/letsencrypt.sh
priority=0
autostart=true
autorestart=false
exitcodes=0
startretries=1
startsecs=1

[program:edg-mkdgridmap-once]
command=/usr/sbin/edg-mkgridmap
    --conf=/etc/edg-mkgridmap.conf
    --output=/etc/grid-security/grid-mapfile
    --safe --cache -quiet
priority=10
autostart=true
autorestart=false
exitcodes=0
startretries=3
startsecs=0

[program:globus-gridftp-server]
command=/usr/sbin/globus-gridftp-server
    -c /etc/gridftp.conf
    -pidfile /var/run/globus-gridftp-server.pid
    -d all
    -l /var/log/gridftp.log
    -no-detach
    -config-base-path /
priority=50
autostart=true
autorestart=true
startretries=3

[program:cron]
command=/usr/sbin/crond -n -s
priority=10
autostart=true
autorestart=true

